# thinkfan/thinkfan.service.d/00-generate-config.conf
[Service]
ExecStartPre=/usr/local/bin/thinkfan-config
