# thinkfan/thinkfan.service.d/10-restart-on-failure.conf
[Unit]
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=3
